# watchdog.timer
[Unit]
Description=Data Collector Watchdog Timer

[Timer]
OnBootSec=120
OnUnitActiveSec=60
AccuracySec=5

[Install]
WantedBy=timers.target
